# data/systemd/facron.service.in
# Copyright 2012-2015 Marc-Antoine Perennou <Marc-Antoine@Perennou.com>
# Distributed under the termes of the GNU Geneal Public License v3

[Unit]
Description=fanotify cron system

[Service]
ExecStart=@sbindir@/facron
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=default.target
